# packaging/remote-prover.service
[Unit]
Description=Remote Aleo Prover
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=remote-prover
Group=remote-prover
EnvironmentFile=-/etc/default/remote-prover
WorkingDirectory=/var/lib/remote-prover
ExecStart=/usr/bin/remote-prover
Restart=on-failure
RestartSec=3

# Resource limits
LimitNOFILE=65536

# Environment
Environment=ALEO_HOME=/var/lib/remote-prover/.aleo

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/remote-prover

[Install]
WantedBy=multi-user.target
